# infrastructure/launchd/ai.hermes.binance-futures-1m-collector.plist.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.hermes.binance-futures-1m-collector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ye/.hermes/scripts/run_binance_futures_1m_collector.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ye/.hermes/scripts</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ye/.hermes/runtime-venv/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ye/.hermes/research/binance_futures_1m_collector_runtime/logs/collector.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ye/.hermes/research/binance_futures_1m_collector_runtime/logs/collector.error.log</string>
</dict>
</plist>
